# launchd/com.immich.accelerator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immich.accelerator</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-m</string>
        <string>accelerator</string>
        <string>watch</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/immich-apple-silicon</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/immich-accelerator.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/immich-accelerator-error.log</string>
</dict>
</plist>
